# tyua-doorbell/rootfs/etc/s6-overlay/s6-rc.d/doorbell/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Example
# Runs example2 script
# ==============================================================================

bashio::log.info "Starting ..."

exec /usr/bin/doorbell.sh